# system/dev/streaming.service
[Service]
Type=simple
User=root
Restart=on-failure
WorkingDirectory=/home/upstage/nms
ExecStart=node app.js

[Install]
WantedBy=multi-user.target

