# config/dev/upstage_event_archive.service
[Unit]
Description=UpStage event archive service

[Service]
User=upstage
Group=upstage
Type=simple
ExecStart=/usr/bin/python3 /home/upstage/upstage/scripts/run_event_archive.py

[Install]
WantedBy=multi-user.target
